# doc/ceilodata2ssm.cron
01 08 * * * root /usr/bin/ceilodata2ssm -l -hst -p -v >> /var/log/cloudaccounting.log 2>&1 
