# doc/cloudaccounting.cron
01 08 * * * root /usr/bin/cloudaccounting -p -v >> /var/log/cloudaccounting.log 2>&1 
