# src/cloudaccounting.cron
01 08 * * * root /usr/libexec/cloudaccounting -p -v >> /var/log/cloudaccounting.log 2>&1 
