# src/cloudaccounting.cron
01 08 * * * root /usr/libexec/cloudaccounting -p 
